# conf/bitcoin_scraper.ini
[program:bitcoin_scraper]
command=python3 scrap_bitcoins.py
directory=/root/how-to-scrapy-bitcoin
process_name=%(program_name)s-%(process_num)02d
numprocs=100
autostart=yes
autorestart=true
redirect_stderr=true
stdout_logfile=/root/how-to-scrapy-bitcoin/test.log
