# conf/bitcoind.ini
[program:bitcoind]
command=/root/bitcoin-0.16.0/bin/bitcoind -rpcuser=ligulfzhou -rpcpassword=zhouligang
process_name=bitcoind
numprocs=1
autostart=yes
autorestart=true
redirect_stderr=true
stdout_logfile=/root/.bitcoin/error.log
