# conf/electrumx.ini
[program:electrumx]
command=/root/electrumx/electrumx_server.py
directory=/root/electrumx
environment=DB_DIRECTORY="/root/.electrumx",DAEMON_URL="http://ligulfzhou:zhouligang@127.0.0.1:8332/",COIN="BitcoinSegwit",ALLOW_ROOT="1",HOST="0.0.0.0",TCP_PORT="50001"
process_name=electrumx
numprocs=1
autostart=yes
autorestart=true
redirect_stderr=true
stdout_logfile=/root/.electrumx/log.log
