# conf/storage_service.ini
[program:storage_service]
command=python3 app.py
directory=/root/how-to-scrapy-bitcoin/storage_service
process_name=%(program_name)s-%(process_num)02d
numprocs=1
autostart=yes
autorestart=true
redirect_stderr=true
stdout_logfile=/root/how-to-scrapy-bitcoin/storage_service.log
